# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2ctl]
command=/usr/sbin/apache2ctl -DFOREGROUND

[program:python3]
command=/usr/bin/python3 /api.py
